# installer/systemd/hal0-agent@hermes.service.d/override.conf
# hermes-specific overrides for the hal0-agent@ template.
#
# Anything under hal0-agent@<id>.service.d/ applies ONLY to that
# instance — `hal0-agent@piccoder.service` would NOT inherit these.
# This keeps the base template generic while letting us pin
# hermes-shaped env without touching the template every release.

[Unit]
Documentation=https://github.com/Hal0ai/hal0/blob/main/docs/agents/hermes/SERVICE.md

[Service]
# HERMES_HOME is what the hermes binary reads to find its config,
# personas, plugins, and skills. hal0 bootstrap claims a marker file
# (`.hal0-managed`) inside this dir; if the marker is missing the
# bootstrap refuses to write here (see
# `src/hal0/agents/hermes_provision._claim_hermes_home`).
Environment="HERMES_HOME=/var/lib/hal0/.hermes"

# Mirror the `--tui` flag via env so any nested hermes-cli invocation
# (e.g. hermes shelling out to itself for a one-shot) sees the same
# dashboard-TUI mode the parent was started in.
Environment="HERMES_DASHBOARD_TUI=1"

# Pin the dashboard's web dist next to the wheel-installed hermes
# binary so `--skip-build` finds it without needing npm at runtime.
# The path is overridable in /etc/hal0/agents/hermes.env if a user
# wants to point at a hand-built tree.
Environment="HERMES_WEB_DIST=/var/lib/hal0/venvs/hermes/lib/python3.12/site-packages/hermes_cli/web_dist"

# Lemonade endpoint hint — Hermes's `hal0` model-provider plugin
# reads this to discover the lemond OpenAI-compatible base. Default
# matches `installer/install.sh`'s lemond bind (127.0.0.1:13305).
Environment="HAL0_LEMONADE_BASE=http://127.0.0.1:13305"
