# installer/systemd/hal0-api.service
[Unit]
Description=hal0 API daemon (FastAPI + dispatcher)
Documentation=https://github.com/hal0-dev/hal0
# Start after the network is up so any upstream health probes can reach
# external providers immediately on boot.
After=network.target

[Service]
# Phase 0: running as root is acceptable for bootstrapping.
# Phase 1 target: change to User=hal0 once the hal0 system user is created
# by install.sh and the FHS layout is in place. The hal0 user needs:
#   - read access to /usr/lib/hal0/current/
#   - read/write access to /var/lib/hal0/
#   - read access to /etc/hal0/
# Uncomment the line below and remove "User=root" when ready:
# User=hal0
User=root

Type=simple

WorkingDirectory=/usr/lib/hal0/current

# /etc/hal0/api.env is written by install.sh (and never clobbered on
# re-runs). It carries HAL0_PORT, HAL0_LOG_LEVEL, and any site-local
# overrides. HAL0_HOME is NOT set here — the FHS layout is the default.
EnvironmentFile=/etc/hal0/api.env

# 'hal0 serve' is the daemon entry point; the binary lives inside the
# versioned venv so the symlink hop is: current/ -> hal0-X.Y.Z/
ExecStart=/usr/lib/hal0/current/bin/hal0 serve

# Restart on any non-zero exit. The 3s back-off matches haloai convention
# and gives the OS a moment to release the port before re-binding.
Restart=on-failure
RestartSec=3

# Feed both journald and the console so a TTY operator sees live output
# without needing 'journalctl -f'.
StandardOutput=journal+console
StandardError=journal+console

# Tag all log lines for easy filtering: journalctl -t hal0-api
SyslogIdentifier=hal0-api

[Install]
WantedBy=multi-user.target
