# installer/systemd/hindsight-api.service
[Unit]
Description=hal0 shared Hindsight memory engine (platform brain)
# lemond provides the extraction LLM (qwen3-it-4b-FLM). Soft ordering: an
# unknown unit name here is ignored by systemd, so this is best-effort.
After=network-online.target hal0-lemonade.service
Wants=network-online.target

[Service]
Type=simple
User=hal0
Group=hal0
WorkingDirectory=/var/lib/hal0/memory/hindsight

# --- pg0 embedded Postgres data root ---
# Hindsight's embedded pg0 lives at $HOME/.pg0 (there is NO HINDSIGHT_API_DATA_DIR
# env var — verified against hindsight-docs 0.7.2). Pin HOME so pg0 persists under
# the hal0 data root rather than the hal0 service user's default home.
Environment=HOME=/var/lib/hal0/memory/hindsight

# --- writable HF cache (spike gotcha #2) ---
# Default HF cache symlinks to read-only /mnt/ai-models and dies downloading the
# local BGE embedder; point HF_HOME at this hal0-writable dir.
Environment=HF_HOME=/var/lib/hal0/memory/hindsight/hf-cache

# --- extraction / consolidation / reflect LLM -> lemond (OpenAI-compatible) ---
# Provider value is `openai` for any OpenAI-compatible endpoint; base_url MUST
# include the /v1 suffix (verified: hindsight-docs models.md examples).
# A NON-EMPTY key is required even for no-auth lemond (spike gotcha #1).
Environment=HINDSIGHT_API_LLM_PROVIDER=openai
Environment=HINDSIGHT_API_LLM_BASE_URL=http://127.0.0.1:13305/v1
Environment=HINDSIGHT_API_LLM_MODEL=qwen3-it-4b-FLM
Environment=HINDSIGHT_API_LLM_API_KEY=lemonade-local-noauth
# NPU extraction is slow (~45-90s/extract, async); raise the LLM timeout so
# fact extraction doesn't abort mid-call (spike model finding).
Environment=HINDSIGHT_API_LLM_TIMEOUT=180

# Embeddings + reranker are LOCAL (bundled BGE 384-d + MiniLM cross-encoder) —
# the full `hindsight-api` package default. No external embed/rerank config needed.
# (P2-5 may repoint the reranker at lemond's bge-reranker-v2-m3 for §4b.)
# CT 105 has no NVIDIA GPU (Strix Halo iGPU only) — force CPU so torch/onnx don't
# probe for CUDA. The bundled CUDA torch wheels are unused dead weight on this box.
Environment=HINDSIGHT_API_EMBEDDINGS_LOCAL_FORCE_CPU=true
Environment=HINDSIGHT_API_RERANKER_LOCAL_FORCE_CPU=true

# Bind to localhost on the pinned port so the provider/client don't chase a
# dynamic port (spike gotcha #3). hal0-api proxies this; not LAN-exposed.
ExecStart=/var/lib/hal0/memory/hindsight/.venv/bin/hindsight-api --host 127.0.0.1 --port 9177

Restart=on-failure
RestartSec=3
TimeoutStartSec=120

[Install]
WantedBy=multi-user.target
