# installer/systemd/hal0-api.service.d/20-run-as-hal0.conf
# hal0-api run-as drop-in (D hardened-perms flip).
#
# Installed ONLY when the installer is run with HAL0_USER != root (the
# hardened-perms flip). It drops hal0-api from root to the unprivileged
# ${HAL0_USER} service account and pins its supplementary primary group to
# the same name, mirroring the hal0-agent@ run-as pattern. The base unit
# already substitutes `User=${HAL0_USER}` inline; this drop-in re-asserts it
# AND adds the matching `Group=` (the generated base unit carries no Group=
# line), so the daemon's group on the shared /etc/hal0 + /var/lib/hal0 setgid
# trees is deterministic.
#
# The `__HAL0_USER__` token below is substituted by installer/install.sh at
# write time (it is NOT a systemd specifier). With HAL0_USER=root this file is
# never installed, so default (root) installs are unchanged.
#
# Prerequisite: the privileged seam (PR #943 — hal0-slotctl + sudoers +
# container.py euid routing) is what lets the dropped daemon still manage slots
# without being root.

[Service]
User=__HAL0_USER__
Group=__HAL0_USER__
