# packaging/systemd/hal0-caddy.service
[Unit]
Description=hal0 Caddy reverse proxy (HTTPS + auth)
Documentation=https://caddyserver.com/docs/
# Network must be up before Caddy binds :443. We don't gate on
# hal0-api.service because Caddy's reverse_proxy handles upstream
# unavailability with a retry; making the API a hard dep would mean a
# brief 502 window during API restarts that today is invisible.
After=network-online.target
Wants=network-online.target
After=hal0-api.service hal0-openwebui.service
# StartLimit lives in [Unit] on modern systemd (Restart= still in [Service]).
StartLimitIntervalSec=120
StartLimitBurst=10

[Service]
Type=notify
# Run as the system 'caddy' user that the Debian/Ubuntu caddy package
# creates. On Arch (where Caddy comes from a release tarball) the
# installer chowns the cert + log dirs to root and the unit runs as root
# — that path is documented in installer/README.md.
User=caddy
Group=caddy
# Caddy needs CAP_NET_BIND_SERVICE to bind 80/443 as a non-root user.
AmbientCapabilities=CAP_NET_BIND_SERVICE
CapabilityBoundingSet=CAP_NET_BIND_SERVICE

# Caddy's data dir holds the auto-mint internal-CA and any ACME state.
# Keeping it under /var/lib/caddy (Debian default) means a hal0 update
# never touches certs.
Environment=XDG_DATA_HOME=/var/lib/caddy
Environment=XDG_CONFIG_HOME=/etc/caddy

# /etc/hal0/Caddyfile is rendered by install.sh from
# packaging/caddy/Caddyfile.template. Re-render with
# `bash installer/install.sh --auth=basic` and `systemctl reload
# hal0-caddy` to apply changes without dropping in-flight connections.
ExecStart=/usr/bin/caddy run --config /etc/hal0/Caddyfile --adapter caddyfile
ExecReload=/usr/bin/caddy reload --config /etc/hal0/Caddyfile --adapter caddyfile --force

# Always restart — a wedged Caddy means the entire dashboard goes dark.
Restart=always
RestartSec=3

# 30s graceful stop window: long enough for in-flight chat-completions
# streams to finish; short enough that systemd doesn't hang on a wedged
# upstream during shutdown.
TimeoutStopSec=30

# Light hardening — Caddy is a network service we don't want to give any
# more host access than it needs.
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
ReadWritePaths=/var/lib/caddy /var/log/caddy

StandardOutput=journal+console
StandardError=journal+console
SyslogIdentifier=hal0-caddy

[Install]
WantedBy=multi-user.target
